# scripts/file-integrity-check.timer
[Unit]
Description=Run FIM check daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
